# systemd/sh-discord-bot.service
[Unit]
Description=SH Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# ВАЖНО: поменяй пути под свою установку
WorkingDirectory=/opt/shbot/SH_discord_bot_split
# Токен задаётся через ENV. Вариант 1 (прямо здесь):
#   Environment=DISCORD_TOKEN=ВАШ_ТОКЕН
# Вариант 2: задать переменную окружения для сервиса через drop-in (recommended).
# (Файл .env проектом больше не используется.)
ExecStart=/opt/shbot/SH_discord_bot_split/.venv/bin/python /opt/shbot/SH_discord_bot_split/main.py
Restart=on-failure
RestartSec=5

# ВАЖНО: лучше запускать под отдельным пользователем
User=shbot
Group=shbot

# Базовая защита
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/shbot/SH_discord_bot_split

[Install]
WantedBy=multi-user.target
